# config/systemd/remind-me-to-move.service
[Unit]
Description="You should try to move regularly, this should remind you"

[Service]
ExecStart=/bin/remind_me_to_move.bash
