# config/autostart/start_alacritty.bash.desktop
[Desktop Entry]
Exec=/bin/start_alacritty.bash
Icon=dialog-scripts
Name=start_alacritty.bash
Path=
Type=Application
X-KDE-AutostartScript=true
