# config/autostart/xswap.bash.desktop
[Desktop Entry]
Exec=/bin/xswap.bash
Icon=dialog-scripts
Name=xswap.bash
Path=
Type=Application
X-KDE-AutostartScript=true
